# load-balancer-servo.cron
@reboot root /bin/bash -c "python /usr/libexec/load-balancer-servo/servo-prep"
