# load-balancer-servo.cron
@reboot root /bin/bash -c "python /usr/libexec/load-balancer-servo/ntp-update"
@reboot root /bin/bash -c "python /usr/libexec/load-balancer-servo/dns-update"
